# backend/supertech.service
[Unit]
Description=SuperTech Django Application
After=network.target

[Service]
Type=notify
User=www-data
WorkingDirectory=/var/www/supertech/backend

Environment="PATH=/var/www/supertech/backend/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=supertech_backend.settings"
Environment="ENVIRONMENT=production"

ExecStart=/var/www/supertech/backend/venv/bin/gunicorn \
    --workers 4 \
    --worker-class sync \
    --bind 127.0.0.1:8000 \
    --timeout 60 \
    --access-logfile /var/log/supertech/access.log \
    --error-logfile /var/log/supertech/error.log \
    supertech_backend.wsgi:application

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
